# supervisor/supervisord-logstash-shipper.conf
[program:logstash-shipper]
command=/opt/services/logstash-shipper.sh
numprocs=1
autostart=true
autorestart=true