# supervisor/supervisord-logstash-indexer.conf
[program:logstash-indexer]
command=/opt/services/logstash-indexer.sh
numprocs=1
autostart=true
autorestart=true