# supervisor/supervisord-elasticsearch.conf
[program:elasticsearch]
command=/opt/services/elasticsearch.sh
numprocs=1
autostart=true
autorestart=true